# laravel-horizon/supervisord.d/laravel-worke-telegram-limited-api-request.conf
[program:TelegramBot-Queue-TelegramLimitedApiRequest]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=TelegramLimitedApiRequest --max-time=3600
autostart=true
autorestart=true
numprocs=2
user=laradock
redirect_stderr=true
